# no-search/supervisord-riak.conf
[supervisord]
nodaemon=true

[program:riak]
command=bash -c "/usr/sbin/riak console"
redirect_stderr=true
stdout_logfile=/var/log/riak/%(program_name)s.log
